# consul-template-haproxy/supervisord-haproxy.conf
[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg
numprocs=1
autostart=true
autorestart=true
