# consul-template-haproxy/supervisord-consul-template.conf
[program:consul-template]
command=consul-template-haproxy
numprocs=1
autostart=true
autorestart=true
