# serf-supervisord/supervisord-serf.ini
[program:serf]
command=/usr/sbin/serf agent -config-dir /etc/serf/serf.d
numprocs=1
autostart=true
autorestart=true
