# consul-supervisord/supervisord-consul.conf
[program:consul]
command=consul-agent
numprocs=1
autostart=true
autorestart=true
