# consul-supervisord/supervisord-consul-join.conf
[program:consul-join]
command=consul-join-master
autorestart=false
startsecs=0
