# serf-aware-haproxy/supervisord-serf-join.ini
[program:serf-join]
command=/serf-join-master.sh
autorestart=false
startsecs=0
