# consul-ruby/supervisord-ruby.conf
[program:consul-template-loader]
command=consul-template-loader
numprocs=1
autostart=true
autorestart=true
